# splunk_app/find_evil/appserver/static/command.a2ui.json
{"version": "v0.9", "createSurface": {"surfaceId": "command", "catalogId": "https://a2ui.org/specification/v0_9/catalogs/basic/catalog.json"}}
{"version": "v0.9", "updateComponents": {"surfaceId": "command", "components": [{"id": "root", "component": "Column", "children": ["title", "verdict", "kpis", "sev_h", "sev", "tech_h", "tech", "lol_h", "lol"]}, {"id": "title", "component": "Text", "variant": "h2", "text": "Forensic Command Center — compromised DC (SRL-2018)"}, {"id": "verdict", "component": "VerdictBadge", "verdict": {"path": "/verdict"}, "summary": {"path": "/summary"}}, {"id": "kpis", "component": "KpiRow", "path": "/kpis"}, {"id": "sev_h", "component": "Text", "variant": "h3", "text": "Detections by severity"}, {"id": "sev", "component": "SeverityBar", "path": "/severity"}, {"id": "tech_h", "component": "Text", "variant": "h3", "text": "Detected MITRE ATT&CK kill-chain"}, {"id": "tech", "component": "TechniqueTable", "path": "/techniques"}, {"id": "lol_h", "component": "Text", "variant": "h3", "text": "Attack tools (LOLBins) in memory"}, {"id": "lol", "component": "LolbinBars", "path": "/lolbins"}]}}
{"version": "v0.9", "updateDataModel": {"surfaceId": "command", "path": "/", "value": {"verdict": "COMPROMISED", "summary": "Windows Server 2016 domain controller — Active Directory credential exfiltration (NTDS.dit + Mimikatz). Agent-driven investigation via the Splunk MCP Server.", "kpis": [{"label": "Critical", "value": 4, "tone": "critical"}, {"label": "High", "value": 7, "tone": "high"}, {"label": "Techniques", "value": 15, "tone": "neutral"}, {"label": "Processes", "value": 122, "tone": "info"}], "severity": [{"sev": "critical", "count": 4}, {"sev": "high", "count": 7}, {"sev": "medium", "count": 2}, {"sev": "informational", "count": 2}], "techniques": [{"signature": "LOLBin_NTDS_Access", "severity": "critical", "mitre": "T1003.003, T1074", "description": "LOLBin tools accessing Active Directory database"}, {"signature": "Credential_Dumping_Framework", "severity": "critical", "mitre": "T1003.001", "description": "Credential dumping framework signatures (mimikatz, etc.)"}, {"signature": "NTDS_Staging_Path", "severity": "critical", "mitre": "T1074.001", "description": "NTDS.dit staged in temp directory"}, {"signature": "NTDS_Extraction_ntdsutil", "severity": "critical", "mitre": "T1003.003", "description": "NTDS.dit extraction via ntdsutil IFM"}, {"signature": "Schtasks_Persistence_Creation", "severity": "high", "mitre": "T1053.005", "description": "Scheduled task creation for persistence"}, {"signature": "LOLBin_Download_Tools", "severity": "high", "mitre": "T1105", "description": "Living-off-the-land download (bitsadmin, certutil, curl)"}, {"signature": "PsExec_Remote_Execution", "severity": "high", "mitre": "T1021.002", "description": "PsExec-style remote execution"}, {"signature": "WMI_Remote_Execution", "severity": "high", "mitre": "T1021.006, T1047", "description": "WMI remote execution for lateral movement"}, {"signature": "PowerShell_Remote_Download", "severity": "high", "mitre": "T1105", "description": "PowerShell downloading remote payload"}, {"signature": "PowerShell_Obfuscated_Execution", "severity": "high", "mitre": "T1059.001, T1027", "description": "PowerShell with obfuscation/encoding"}, {"signature": "Shadow_Copy_Credential_Theft", "severity": "high", "mitre": "T1003, T1006", "description": "Credential extraction via VSS shadow copies"}, {"signature": "PowerShell_Remoting_Session", "severity": "medium", "mitre": "T1021.006", "description": "PowerShell Remoting (WinRM/WS-Man) session"}, {"signature": "Recycle_Bin_Antiforensic", "severity": "medium", "mitre": "T1070", "description": "Move to Recycle Bin for anti-forensic cleanup"}, {"signature": "Defender_Blocked_Execution", "severity": "informational", "mitre": "", "description": "Windows Defender blocked malicious execution"}, {"signature": "SRL2018_Scenario_Markers", "severity": "informational", "mitre": "", "description": "Specific markers for SRL-2018 scenario (shieldbase.lan domain)"}], "lolbins": [{"name": "cmd.exe", "count": 25}, {"name": "powershell.exe", "count": 1}]}}}