# splunk_app/find_evil/appserver/static/forensic_report.a2ui.json
{"version": "v0.9", "createSurface": {"surfaceId": "forensic", "catalogId": "https://a2ui.org/specification/v0_9/catalogs/basic/catalog.json"}}
{"version": "v0.9", "updateComponents": {"surfaceId": "forensic", "components": [{"id": "root", "component": "Column", "children": ["title", "verdict", "kpis", "sev_h", "sev", "tech_h", "tech", "reco_h", "reco", "analysis"]}, {"id": "title", "component": "Text", "variant": "h2", "text": "Incident report — splunklib.ai agent → A2UI"}, {"id": "verdict", "component": "VerdictBadge", "verdict": {"path": "/verdict"}, "summary": {"path": "/summary"}}, {"id": "kpis", "component": "KpiRow", "path": "/kpis"}, {"id": "sev_h", "component": "Text", "variant": "h3", "text": "Severity breakdown"}, {"id": "sev", "component": "SeverityBar", "path": "/severity"}, {"id": "tech_h", "component": "Text", "variant": "h3", "text": "MITRE ATT&CK kill-chain"}, {"id": "tech", "component": "TechniqueTable", "path": "/techniques"}, {"id": "reco_h", "component": "Text", "variant": "h3", "text": "Remediation recommendations"}, {"id": "reco", "component": "RecommendationList", "path": "/recommendations"}, {"id": "analysis", "component": "Collapsible", "title": "Detailed agent analysis (kill-chain & reasoning)", "text": {"path": "/analysis"}}]}}
{"version": "v0.9", "updateDataModel": {"surfaceId": "forensic", "path": "/", "value": {"verdict": "COMPROMISED", "summary": "4 critical detections, 7 high — 15 MITRE techniques on the domain controller (122 processes analyzed).", "analysis": "## Incident Analysis — Memory Forensics Report\n\n### Verdict: ⛔ COMPROMISED\n\nThe memory image yields **15 YARA hits** distributed across **4 critical**, **7 high**, **2 medium**, and **2 informational** severity levels. The triage summary confirms **13 distinct MITRE ATT&CK technique groups** are represented. Together, the detections describe a **full-spectrum, hands-on-keyboard intrusion** against a Windows Active Directory environment (`shieldbase.lan`), covering every phase of the kill-chain from initial access through credential theft, lateral movement, persistence, and active anti-forensics.\n\n---\n\n### Kill-Chain Reconstruction\n\n#### Stage 1 — Initial Access & Execution\n- **`LOLBin_Download_Tools`** (T1105) and **`PowerShell_Remote_Download`** (T1105): The attacker leveraged built-in Windows binaries (`bitsadmin`, `certutil`, `curl`) and PowerShell to download remote payloads, likely bypassing network-layer controls that block known offensive tooling. These two signatures are the only ones with **2 detections each**, underlining repeated use.\n- **`PowerShell_Obfuscated_Execution`** (T1059.001 / T1027): Encoded or otherwise obfuscated PowerShell was executed to evade signature-based detection. Windows Defender flagged at least one execution attempt (`Defender_Blocked_Execution`), but the attacker succeeded via alternative vectors.\n\n#### Stage 2 — Lateral Movement\n- **`PsExec_Remote_Execution`** (T1021.002): PsExec-style remote execution over SMB was used to move between hosts within the domain.\n- **`WMI_Remote_Execution`** (T1021.006 / T1047): WMI was weaponized for remote command execution, a technique that frequently blends with legitimate administrative traffic.\n- **`PowerShell_Remoting_Session`** (T1021.006): WinRM/WS-Man sessions confirm additional PowerShell-based lateral movement channels.\n\n#### Stage 3 — Credential Access *(Most Critical Phase)*\nFour **critical-severity** signatures converge on a systematic, multi-vector Active Directory credential harvesting operation:\n- **`Credential_Dumping_Framework`** (T1003.001): Mimikatz or equivalent framework artifacts detected in memory — LSASS process memory was targeted for plaintext credentials and NT hashes.\n- **`Shadow_Copy_Credential_Theft`** (T1003 / T1006): Volume Shadow Copies (VSS) were accessed to extract SAM/SYSTEM/NTDS credential material offline, circumventing LSASS-level protections.\n- **`NTDS_Extraction_ntdsutil`** (T1003.003): The `ntdsutil` IFM (Install From Media) feature was abused to generate a full offline copy of `NTDS.dit` — the Active Directory database containing password hashes for **every domain account**.\n- **`LOLBin_NTDS_Access`** (T1003.003 / T1074): LOLBin utilities were used to directly access the AD database, reinforcing the ntdsutil finding.\n- **`NTDS_Staging_Path`** (T1074.001): The extracted `NTDS.dit` was staged in a temporary directory, strongly suggesting imminent or completed exfiltration.\n\n> ⚠️ **Critical Impact:** A complete NTDS.dit extraction means ALL domain account credentials — including the KRBTGT hash — must be treated as fully compromised. This enables offline password cracking and Golden Ticket attacks with indefinite domain persistence.\n\n#### Stage 4 — Persistence\n- **`Schtasks_Persistence_Creation`** (T1053.005): A scheduled task was created to ensure the attacker's foothold survives reboots, user logoffs, and initial remediation attempts.\n\n#### Stage 5 — Defense Evasion & Anti-Forensics\n- **`Recycle_Bin_Antiforensic`** (T1070): Attacker tooling or evidence files were moved to the Recycle Bin to impede forensic analysis, confirming deliberate post-exploitation cleanup.\n\n#### Contextual Indicator\n- **`SRL2018_Scenario_Markers`**: Environment-specific markers for the `shieldbase.lan` domain confirm the target identity. Treated as contextual data only.\n\n---\n\n### Severity Breakdown (from Triage Summary)\n| Severity | Count |\n|---|---|\n| Critical | 4 |\n| High | 7 |\n| Medium | 2 |\n| Informational | 2 |\n| **Total** | **15** |\n\n---\n\n### Conclusion\nThis is a **confirmed, advanced, multi-stage intrusion**. The attacker achieved full domain-level credential access via three independent methods (LSASS dump, VSS, ntdsutil), established persistence, and actively attempted to destroy forensic evidence. The entire `shieldbase.lan` Active Directory domain must be treated as fully compromised.", "kpis": [{"label": "Critical", "value": 4, "tone": "critical"}, {"label": "High", "value": 7, "tone": "high"}, {"label": "Techniques", "value": 15, "tone": "neutral"}, {"label": "Processes", "value": 122, "tone": "info"}], "severity": [{"sev": "critical", "count": 4}, {"sev": "high", "count": 7}, {"sev": "medium", "count": 2}, {"sev": "informational", "count": 2}], "techniques": [{"signature": "LOLBin_NTDS_Access", "severity": "critical", "mitre": "T1003.003, T1074", "description": "LOLBin tools accessing Active Directory database"}, {"signature": "Credential_Dumping_Framework", "severity": "critical", "mitre": "T1003.001", "description": "Credential dumping framework signatures (mimikatz, etc.)"}, {"signature": "NTDS_Staging_Path", "severity": "critical", "mitre": "T1074.001", "description": "NTDS.dit staged in temp directory"}, {"signature": "NTDS_Extraction_ntdsutil", "severity": "critical", "mitre": "T1003.003", "description": "NTDS.dit extraction via ntdsutil IFM"}, {"signature": "Schtasks_Persistence_Creation", "severity": "high", "mitre": "T1053.005", "description": "Scheduled task creation for persistence"}, {"signature": "LOLBin_Download_Tools", "severity": "high", "mitre": "T1105", "description": "Living-off-the-land download (bitsadmin, certutil, curl)"}, {"signature": "PsExec_Remote_Execution", "severity": "high", "mitre": "T1021.002", "description": "PsExec-style remote execution"}, {"signature": "WMI_Remote_Execution", "severity": "high", "mitre": "T1021.006, T1047", "description": "WMI remote execution for lateral movement"}, {"signature": "PowerShell_Remote_Download", "severity": "high", "mitre": "T1105", "description": "PowerShell downloading remote payload"}, {"signature": "PowerShell_Obfuscated_Execution", "severity": "high", "mitre": "T1059.001, T1027", "description": "PowerShell with obfuscation/encoding"}, {"signature": "Shadow_Copy_Credential_Theft", "severity": "high", "mitre": "T1003, T1006", "description": "Credential extraction via VSS shadow copies"}, {"signature": "PowerShell_Remoting_Session", "severity": "medium", "mitre": "T1021.006", "description": "PowerShell Remoting (WinRM/WS-Man) session"}, {"signature": "Recycle_Bin_Antiforensic", "severity": "medium", "mitre": "T1070", "description": "Move to Recycle Bin for anti-forensic cleanup"}, {"signature": "Defender_Blocked_Execution", "severity": "informational", "mitre": "", "description": "Windows Defender blocked malicious execution"}, {"signature": "SRL2018_Scenario_Markers", "severity": "informational", "mitre": "", "description": "Specific markers for SRL-2018 scenario (shieldbase.lan domain)"}], "recommendations": [{"priority": "•", "horizon": null, "text": "[P1 — IMMEDIATE] Isolate the compromised host and preserve forensic evidence. Take the system offline immediately to halt any ongoing lateral movement or exfiltration channel. Preserve the memory image and full disk image under chain-of-custody before any remediation. Do not simply reimage — the staged NTDS.dit path, scheduled task details, and attacker tooling must be fully documented for scoping."}, {"priority": "•", "horizon": null, "text": "[P1 — IMMEDIATE] Execute a full Active Directory credential reset, starting with KRBTGT (twice). Because NTDS.dit was fully extracted and staged, every domain account hash is compromised. Reset the KRBTGT account password twice in sequence (required to invalidate all existing Kerberos tickets and eliminate the Golden Ticket attack surface), then reset all Tier-0 and Tier-1 privileged accounts (Domain Admins, Enterprise Admins, service accounts), followed by all remaining user accounts. Coordinate with business stakeholders for a planned reset window."}, {"priority": "•", "horizon": null, "text": "[P1 — IMMEDIATE] Hunt for persistence and lateral spread across the entire domain. The scheduled task, PsExec, and WMI lateral movement detections imply additional hosts may be compromised. Immediately enumerate all recently created scheduled tasks, WMI event subscriptions, new local/domain accounts, and modified registry run keys across all domain-joined machines. Treat any system the attacker could have reached via SMB or WMI as potentially compromised until proven otherwise."}, {"priority": "•", "horizon": null, "text": "[P2 — SHORT-TERM] Harden credential access attack surfaces. Enable and enforce Credential Guard (blocks LSASS memory access by hypervisor isolation) and add all privileged accounts to the Protected Users security group (disables NTLM, DES, RC4 for those accounts). Enable LSA Protection (RunAsPPL) to block unauthorized LSASS access. Restrict `ntdsutil`, `vssadmin`, and direct NTDS.dit file access via AppLocker or Windows Defender Application Control (WDAC) policies."}, {"priority": "•", "horizon": null, "text": "[P3 — MEDIUM-TERM] Strengthen detection and reduce the LOLBin/PowerShell attack surface. Deploy PowerShell Script Block Logging and Module Logging across all endpoints to capture obfuscated execution in future incidents. Block or alert on LOLBin download proxies (`certutil -urlcache`, `bitsadmin /transfer`, `curl`) via EDR policy or WDAC rules. Tune EDR/SIEM to alert on `ntdsutil IFM`, `vssadmin create shadow`, and any process opening a handle to `ntds.dit` outside of the NTDS service."}]}}}