# splunk_app/find_evil/appserver/static/incidents.a2ui.json
{"version": "v0.9", "createSurface": {"surfaceId": "incidents", "catalogId": "https://a2ui.org/specification/v0_9/catalogs/basic/catalog.json"}}
{"version": "v0.9", "updateComponents": {"surfaceId": "incidents", "components": [{"id": "root", "component": "Column", "children": ["title", "list"]}, {"id": "title", "component": "Text", "variant": "h2", "text": "SOC incidents — automated AI triage"}, {"id": "list", "component": "IncidentList", "path": "/incidents"}]}}
{"version": "v0.9", "updateDataModel": {"surfaceId": "incidents", "path": "/", "value": {"incidents": [{"time": "2026-06-12T18:29:44.000+02:00", "verdict": "COMPROMISED", "critical": 4, "high": 7, "analysis": "# SOC Analysis Report\n\n## 1) VERDICT: **COMPROMISED**\n\n## 2) MITRE ATT&CK Kill-Chain Summary:\n\n**Initial Access → Execution → Persistence → Credential Access → Collection → Staging → Defense Evasion**\n\n- **Lateral Movement (T1021)**: PsExec, WMI, and PowerShell remoting used to move across the domain\n- **Execution (T1059.001, T1047)**: Obfuscated PowerShell scripts and WMI commands executed\n- **Persistence (T1053.005)**: Scheduled tasks created for maintaining access\n- **Credential Access (T1003.001, T1003.003)**: \n  - Active credential dumping frameworks (Mimikatz-like tools)\n  - NTDS.dit extraction via ntdsutil and shadow copies\n  - Direct NTDS database access via LOLBins\n- **Collection & Staging (T1074.001)**: NTDS database staged for exfiltration\n- **Defense Evasion (T1027, T1070)**: Obfuscation techniques and anti-forensic activities (Recycle Bin manipulation)\n\n## 3) Prioritized Remediation Actions:\n\n### **IMMEDIATE (Next 15 minutes):**\n1. **Isolate the Domain Controller** - Disconnect from network immediately while maintaining logging. Assume all domain credentials are compromised. Activate backup DC for authentication services.\n\n### **URGENT (Next 1 hour):**\n2. **Reset All Privileged Credentials** - Force password resets for all Domain Admins, Enterprise Admins, service accounts, and krbtgt account (twice, 24hrs apart). Revoke all active Kerberos tickets.\n\n### **HIGH PRIORITY (Next 4 hours):**\n3. **Hunt for Lateral Movement** - Investigate all systems accessed via PsExec/WMI/PSRemoting from this DC. Search for staged NTDS files on network shares, identify exfiltration attempts, and contain affected endpoints before credential resets.\n\n**Note**: Engage incident response team and consider full domain rebuild if attacker dwell time is extensive."}]}}}